# deploy/clinicfocus.service
[Unit]
Description=ClinicFocus FastAPI Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/clinic-focus
Environment=AWS_REGION=us-east-1
ExecStart=/home/ec2-user/clinic-focus/.venv/bin/uvicorn be.handlers.api:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
